# systemd/to-drive-uploader.service
[Unit]
Description=Telegram to Google Drive Uploader Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/python bot.py
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

# Keep logs in journald.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
